# openwrt/nucleus/feeds/application/netopeer/files/nc_conn_dbg_keeper.initd
#!/bin/sh /etc/rc.common

STOP=91
NC_CONN_DBG_TO_FLASH=/usr/sbin/nc_conn_dbg_to_flash.sh

stop(){
	${NC_CONN_DBG_TO_FLASH} connection
	${NC_CONN_DBG_TO_FLASH} traceroute
	${NC_CONN_DBG_TO_FLASH} tcpping
}
