# openwrt/overlay/ramips/base-files/files/etc/init.d/ubi
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=00
start() {
	awk '/^[^#]/ { system("ubiattach -p "$1" -d "$2) }' /etc/ubitab
}
