# openwrt/overlay/ramips/base-files/files/etc/init.d/warm_start
#!/bin/sh /etc/rc.common
# Copyright (C) 2013-2014 OpenWrt.org

STOP=90

stop() {
	touch /etc/zyxel/ftp/keep/warm_start
}
