# openwrt/overlay/ramips/base-files/files/etc/init.d/zy_fwversion
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=01

boot() {
	install -D /rom/etc/fwversion /rw/fwversion
}
